# crontab
0 * * * * /home/furaii/backup_databases.sh
